# gunicorn.service
# Systemd service file for Gunicorn
# Salva questo file come: /etc/systemd/system/gunicorn.service

[Unit]
Description=Gunicorn daemon for Portfolio Django Project
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/percorso/alla/cartella/html/portfolio_project
Environment="PATH=/percorso/alla/cartella/html/portfolio_project/venv/bin"
EnvironmentFile=/percorso/alla/cartella/html/portfolio_project/.env
ExecStart=/percorso/alla/cartella/html/portfolio_project/venv/bin/gunicorn \
    --config /percorso/alla/cartella/html/portfolio_project/gunicorn_config.py \
    portfolio_project.wsgi:application

# Restart policy
Restart=always
RestartSec=3

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
